# packaging/rpm/systemd/infrasense-frontend.service
[Unit]
Description=InfraSense Frontend (Nginx)
Documentation=https://github.com/infrasense/infrasense-platform
After=network.target network-online.target infrasense-api.service
Wants=network-online.target

[Service]
Type=forking
ExecStartPre=/usr/sbin/nginx -t -c /etc/infrasense/nginx.conf
ExecStart=/usr/sbin/nginx -c /etc/infrasense/nginx.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
PIDFile=/run/infrasense-nginx.pid
Restart=on-failure
RestartSec=10s
TimeoutStartSec=30s
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=infrasense-frontend

[Install]
WantedBy=multi-user.target
