# packaging/deb/etc/systemd/system/infrasense-api.service
[Unit]
Description=InfraSense API Server
After=network.target postgresql.service

[Service]
Type=simple
User=infrasense
Group=infrasense
ExecStart=/usr/local/bin/infrasense-api
EnvironmentFile=-/etc/infrasense/infrasense-api.env
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
